# src/deb/control/postinst
#!/bin/sh

if which invoke-rc.d >/dev/null 2>&1; then
	invoke-rc.d oss-cluster start
else
	/etc/init.d/oss-cluster start
fi
